# fantastic_vuln/base/lamp/supervisor_config/conf.d/mysql.conf
[program:mysql]
directory=/
autostart=false # should be start manually
autorestart=true
startretries=5
priority=4
command=mysqld
stdout_logfile=/var/log/mysql.log
redirect_stderr=true
